# systemd/firewall-tui.service.example
[Unit]
Description=Personal Firewall TUI (example)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/personal-firewall/tui/firewall_tui.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
